# client-python/systemd/scoreboard.service
[Unit]
Description=LED Matrix Scoreboard for Scout Patrols
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/OsmDeviceAdapter/client-python/src

# Environment variables (update these paths and values)
Environment="API_BASE_URL=https://your-server.example.com"
Environment="CLIENT_ID=scoreboard-rpi"
Environment="POLL_INTERVAL=30"
Environment="TOKEN_FILE=/var/lib/scoreboard/token.txt"
Environment="LOG_LEVEL=INFO"
Environment="SIMULATE_DISPLAY=false"

# Run the scoreboard application
ExecStart=/usr/bin/python3 /home/pi/OsmDeviceAdapter/client-python/src/scoreboard.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scoreboard

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
